# launchd/ai.zai-monitor.alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.zai-monitor.alerts</string>

    <key>ProgramArguments</key>
    <array>
        <!-- TODO: replace USERNAME and the project path with your own -->
        <string>/Users/USERNAME/dev/zai-monitor/.venv/bin/python</string>
        <string>/Users/USERNAME/dev/zai-monitor/alerts.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/dev/zai-monitor</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/zai-monitor-alerts.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/zai-monitor-alerts.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
